# systemd/.config/systemd/user/battery.service
[Unit]
Description=Battery status 

[Service]
ExecStart=/home/sven/.scripts/battery.sh

[Install]
WantedBy=multi-user.target
