# systemd/.config/systemd/user/battery.timer
[Unit]
Description=Check battery every minute.

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
